# attack-chain-outliner-workspace/iteration-1/linux-cron-persistence/eval_metadata.json
{
  "eval_id": 1,
  "eval_name": "linux-cron-persistence",
  "prompt": "Write up an attack chain document for cron-based persistence on Linux servers. Use Markdown format.",
  "assertions": [
    {
      "name": "has_all_8_sections",
      "description": "Document contains all 8 required section headings: Technique Overview, Attack Steps, Detection Opportunities, Data Sources & Log Fields, Pseudocode Queries, and MITRE ATT&CK Mapping"
    },
    {
      "name": "correct_mitre_id",
      "description": "References T1053.003 specifically (Scheduled Task/Job: Cron)"
    },
    {
      "name": "fidelity_ratings_present",
      "description": "At least one detection method is labeled with High, Medium, or Low fidelity"
    },
    {
      "name": "linux_datasources_present",
      "description": "Data sources section mentions auditd or journald or syslog as Linux-specific log sources"
    },
    {
      "name": "pseudocode_uses_udm_fields",
      "description": "Pseudocode queries use UDM-style field names like event.type, process.name, file.path rather than vendor-specific syntax"
    }
  ]
}
