# attack-chain-outliner-workspace/iteration-1/linux-cron-persistence/with_skill/grading.json
{
  "expectations": [
    {
      "text": "Document contains all 8 required section headings: Technique Overview, Attack Steps, Detection Opportunities, Data Sources & Log Fields, Pseudocode Queries, and MITRE ATT&CK Mapping",
      "passed": true,
      "evidence": "All six named sections are present as H3 headings: 'Technique Overview' (line 5), 'Attack Steps' (line 23), 'Detection Opportunities' (line 78), 'Data Sources & Log Fields' (line 104), 'Pseudocode Queries' (line 193), 'MITRE ATT&CK Mapping' (line 267)."
    },
    {
      "text": "References T1053.003 specifically (Scheduled Task/Job: Cron)",
      "passed": true,
      "evidence": "MITRE ATT&CK Mapping table at lines 271-273 explicitly lists 'T1053.003 | Scheduled Task/Job: Cron' across three tactic rows (Persistence, Execution, Privilege Escalation)."
    },
    {
      "text": "At least one detection method is labeled with High, Medium, or Low fidelity",
      "passed": true,
      "evidence": "Detection Opportunities section contains multiple fidelity-labeled entries: 'Crontab File Creation/Modification (High Fidelity)' (line 80), 'User Crontab Modification via crontab Command (Medium Fidelity)' (line 82), 'crond Spawning Unusual Child Processes (High Fidelity)' (line 84). Pseudocode queries also use fidelity headings (lines 195, 202, 209, 217, 224, 232, 239)."
    },
    {
      "text": "Data sources section mentions auditd or journald or syslog as Linux-specific log sources",
      "passed": true,
      "evidence": "Data Sources section explicitly covers 'auditd (Linux Audit Framework)' (line 106) with syscall-level field details, and 'syslog / journald (cron execution logs)' (line 133) with example journalctl commands and syslog paths for both Debian/Ubuntu and RHEL/CentOS."
    },
    {
      "text": "Pseudocode queries use UDM-style field names like event.type, process.name, file.path rather than vendor-specific syntax",
      "passed": true,
      "evidence": "All pseudocode queries use UDM-style dotted field names: 'event.type', 'file.path', 'process.name', 'process.parent.name', 'process.command_line', 'user.name', 'network.destination.ip', 'network.destination.port', 'network.direction', 'network.protocol'. No vendor-specific syntax (SPL, KQL, Sigma) is used."
    }
  ],
  "summary": {
    "passed": 5,
    "failed": 0,
    "total": 5,
    "pass_rate": 1.0
  },
  "claims": [],
  "eval_feedback": {
    "suggestions": [
      "Assertion 1 references '8 required section headings' but names only 6 — the count discrepancy should be resolved in the assertion wording.",
      "Consider adding Low fidelity labeled detections for completeness of the fidelity spectrum."
    ],
    "overall": "Excellent output. All assertions pass. The document follows the required schema precisely, uses UDM-style pseudocode consistently, applies fidelity labels throughout, and covers Linux-specific log sources (auditd, journald, syslog) in appropriate detail."
  }
}
