# kql-mde-xdr-workspace/iteration-1/eval-macos-persistence/with_skill/outputs/query.kql
// MITRE ATT&CK: T1543.001 - Create or Modify System Process: Launch Agent/Launch Daemon
// Detects suspicious LaunchAgent/LaunchDaemon plist modifications on macOS devices

let macOSDevices =
    DeviceInfo
    | where OSPlatform == "macOS"
    | distinct DeviceId;
DeviceEvents
| where Timestamp > ago(7d)
| where DeviceId in (macOSDevices)
| where ActionType == "PlistPropertyModified"
| extend PlistDetails = parse_json(AdditionalFields)
| where PlistDetails.FileName has_any (
    "LaunchAgents",
    "LaunchDaemons"
)
| project
    Timestamp,
    DeviceName,
    PlistPath = tostring(PlistDetails.FileName),
    PropertyName = tostring(PlistDetails.PropertyName),
    PropertyValue = tostring(PlistDetails.PropertyValue),
    InitiatingProcessFileName,
    InitiatingProcessFolderPath,
    InitiatingProcessCommandLine,
    InitiatingProcessAccountName,
    InitiatingProcessSHA256
| order by Timestamp desc
